# RESOURCES/OPTIONS/ssh/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

if [ -f /app/app ]; then
    chmod 755 /app/app
    /app/app&
fi

/opt/etc/init.d/rc.unslung start
exit 0
